# systemd/mvm-status.service
[Unit]
Description=Midnight Validator Monitor - Status Check
After=network-online.target

[Service]
Type=oneshot
User=%USER%
WorkingDirectory=/opt/midnight/mvm
Environment="MVM_DB_PATH=/opt/midnight/mvm/data/mvm.db"
ExecStart=/opt/midnight/mvm/bin/mvm status --once
StandardOutput=journal
StandardError=journal
